# service/debian_client.service
[Unit]
Description=serverstatus clinet service
After=network.target

[Service]
Type=simple
WorkingDirectory=${ClientPath}
ExecStart=python3 ${ClientFile}

[Install]
WantedBy=multi-user.target